# config_files/restic/restic.timer
[Unit]
Decription=Timer for Restic backup service

[Timer]
OnCalendar=*-*-* 23:45:00
Persistance=yes

[Install]
WantedBy=timers.target